# recipes-core/systemd/files/fluke-cda-caldera/launchApp.service
[Unit]
After=instrument.service

Description=Autostart the application's launch script

[Service]
Type=oneshot
ExecStart=/home/Caldera/bin/caldera
RemainAfterExit=yes

[Install]

WantedBy=multi-user.target
